# packages/mylife-home-core-go/mylife-home-core.initd
#!/sbin/openrc-run

name=mylife-home-core
command=/usr/bin/mylife-home-core
command_args=--config=/etc/mylife-home/config.yaml
output_log=/var/log/mylife-home/output.log
error_log=/var/log/mylife-home/error.log

supervisor=supervise-daemon
# respawn_delay=1 # default = immediately
respawn_max=0

depend() {
  need net
  after firewall
}
