# packages/mylife-home-core/mylife-home-core.initd
#!/sbin/openrc-run

name=mylife-home-core
command=/usr/bin/node
command_args=/usr/mylife-home/bin.js
supervisor="supervise-daemon"

depend() {
  need net
  after firewall
}
