# packages/mylife-home-core/mylife-home-core.initd
#!/sbin/openrc-run

name=mylife-home-core
command="node /usr/mylife-home/bin.js"
supervisor="supervise-daemon"

depend() {
  need net
  after firewall
}
